# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
user=root

[program:tts-service]
command=python3 /app/tts-service/index.py
directory=/app/tts-service
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/tts-service.log
stderr_logfile=/var/log/supervisor/tts-service-error.log
environment=PYTHONIOENCODING="utf-8"

[program:nextjs]
command=bun server.js
directory=/app
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/nextjs.log
stderr_logfile=/var/log/supervisor/nextjs-error.log
environment=NODE_ENV="production",PORT="3000"
